# pr2/bin/test2.service
[Unit]
Description=test2
After=network.target
 
[Service]
Type=simple
User=admin
ExecStart=/bin/bash /opt/test2/test2.sh
Restart=on-failure
RestartSec=5
 
[Install]
WantedBy=multi-user.target