# pr1/bin/test1.service
[Unit]
Description=test1
After=network.target
 
[Service]
Type=simple
User=admin
ExecStart=/bin/bash /opt/test1/test1.sh
Restart=on-failure
RestartSec=5
 
[Install]
WantedBy=multi-user.target